# omniverse/Galaxea R1/navigation/galaxea_nav.service
[Unit]
Description=Galaxea R1 Navigation ROS2 Bridge Node
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/galaxea_nav_combined_repo/galaxea_nav_ws
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source install/setup.bash && ros2 launch ros_bridge bridge.launch.py'
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target